# raspberry/etc/udev/rules.d/99-gpib.rules
KERNEL=="gpib[0-9]*", MODE="0660", GROUP="gpib",  RUN+="/usr/local/sbin/gpib_config --minor 0"

